# rootfs-overlay/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=/sbin/agetty --autologin brlpi --noclear %I $TERM
